# usb-reset.service
[Unit]
Description=Power-cycle touchscreen USB port after desktop starts
After=display-manager.service udev.service
Wants=display-manager.service

[Service]
Type=oneshot
# Tweak timing and target here without editing the script
Environment=SLEEP_BEFORE=10
Environment=HUB=1-1
Environment=PORT=1
ExecStart=/usr/local/bin/usb_reset.sh

[Install]
WantedBy=graphical.target
