# config/autostart/kiosk.desktop
[Desktop Entry]
Type=Application
Name=Kiosk
Exec=__REPO_ROOT__/scripts/kiosk.sh
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
